# deploy/systemd/sanuvia.service
# Sanuvia AI — uvicorn service.  Install to /etc/systemd/system/sanuvia.service
# Assumes the repo is cloned at /opt/sanuvia/my-voice-ai and owned by user `sanuvia`.
[Unit]
Description=Sanuvia AI (FastAPI + uvicorn)
After=network.target

[Service]
Type=simple
User=sanuvia
Group=sanuvia
WorkingDirectory=/opt/sanuvia/my-voice-ai

# keep model caches inside the app dir (writable by the service user)
Environment=HF_HOME=/opt/sanuvia/my-voice-ai/.cache/huggingface
Environment=HF_HUB_DISABLE_XET=1
# low-memory friendly: cap threads and glibc malloc arenas (keeps RSS ~0.6 GB).
# Raise OMP/SHERPA threads to your vCPU count if the box has spare CPU + RAM.
Environment=OMP_NUM_THREADS=2
Environment=SHERPA_NUM_THREADS=2
Environment=MALLOC_ARENA_MAX=2

# AI layer (assistant / supervisor / memory). Disabled until a Groq key is set.
# NOTE: enabling sends transcript text to Groq — gate behind consent for clinical use.
# Environment=GROQ_API_KEY=gsk_xxxxxxxxxxxxxxxxxxxx

ExecStart=/opt/sanuvia/my-voice-ai/.venv/bin/uvicorn backend.main:app --host 127.0.0.1 --port 8000

Restart=on-failure
RestartSec=5
# model loading on first boot can take a while
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
